# deploy/pronote-school.service
[Unit]
Description=Pronote Automation — school-hours sync (vie scolaire)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pronote
WorkingDirectory=/opt/pronote
ExecStart=/usr/bin/java \
    -Xmx128m \
    -Xms32m \
    -DLOG_LEVEL=INFO \
    -jar /opt/pronote/pronote-automation.jar \
    --config /opt/pronote/config.yaml \
    --features schoolLife
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pronote-school

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/pronote/data
